# rel/files/upstart/ddb.conf
description "Dalmatiner DB"
start on filesystem
stop on runlevel [06]
exec /usr/lib/ddb/bin/ddb console -noinput +Bd
